# novosga-2.2-standalone/config/crontab
# m h  dom mon dow   command
0 1 * * * /var/www/html/bin/console novosga:reset
