# novosga-2.1/supervisor/apache2.conf
[program:apache2]
command = apache2-foreground
stdout_logfile = /var/log/novosga-apache2.log
stderr_logfile = /var/log/novosga-apache2-error.log
