# novosga-2.1/supervisor/websocket.conf
[program:websocket]
command = /var/www/html/bin/console novosga:websocket start
stdout_logfile = /var/log/novosga-websocket.log
stderr_logfile = /var/log/novosga-websocket-error.log
